# pkg/blacksharkd.service
[Unit]
Description=Razer BlackShark V3 Pro headset daemon
Documentation=https://github.com/yourusername/blackshark-linux
# Restart if the headset is unplugged and replugged — udev can trigger this later
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/blacksharkd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
